# lab9/pon.service
[Unit]
Description=PONPONPON
After=sshd.service

[Service]
ExecStart=/usr/bin/printf "PONPONPON"

[Install]
WantedBy=multi-user.target